# kindle/crontab
0 * * * * /mnt/us/weather/display-weather.sh
